# onyx.service
# /etc/systemd/system/onyx.service

[Unit]
Description=Onyx Security Orchestrator
Documentation=https://github.com/DiscoMouse/onyx
After=network.target network-online.target
Requires=network-online.target

[Service]
Type=notify
User=onyx
Group=onyx

# PATHS - Updated to use the 'proxy' subcommand
ExecStart=/usr/bin/onyx proxy run --environ --config /etc/onyx/Caddyfile
ExecReload=/usr/bin/onyx proxy reload --config /etc/onyx/Caddyfile --force

# ENVIRONMENT CONFIGURATION
Environment="ONYX_ENV=dev"

# SECRETS MANAGEMENT
EnvironmentFile=-/etc/onyx/onyx.env

# PROCESS LIMITS
TimeoutStopSec=5s
LimitNOFILE=1048576
LimitNPROC=512

# HARDENING
PrivateTmp=true
ProtectSystem=full
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target