# onyx.service
[Unit]
Description=Onyx Security Orchestrator
Documentation=https://github.com/DiscoMouse/onyx
After=network.target network-online.target
Requires=network-online.target

[Service]
# 'simple' is best for Go apps until we implement systemd notification logic
Type=simple
User=onyx
Group=onyx

# CLEAN STARTUP: start the binary
ExecStart=/usr/bin/onyx

# RESTART POLICY: If it crashes, bring it back
Restart=on-failure
RestartSec=5s

# ENVIRONMENT
Environment="ONYX_ENV=production"

# CAPABILITIES: Allows non-root bind to the ports
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# HARDENING
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target